# weather.service
[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=etc/default/weather
ExecStart=/sbin/weather.sh

[Install]
WantedBy=multi-user.target
